# systemd/.config/systemd/user/hintsd.service
[Unit]
Description=Hints daemon
[Service]
Type=simple
ExecStart=/home/kike/.local/bin/hintsd
Restart=always
[Install]
WantedBy=default.target
